# confs/celery.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn
setuid ubuntu
setgid ubuntu

script
  /bin/bash /home/ubuntu/apps/politico-civic/run_on_server.sh celery worker --app=civic.celery.app -l info -f {{ SERVER_LOG_PATH }}/celery.log
end script